# systemd_units/phantom-qa.service
[Unit]
Description=Meta App Factory Phantom QA Daemon (FastAPI)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/meta_app_factory/backend/Phantom_QA_Elite/backend
Environment="PATH=/var/www/meta_app_factory/backend/venv/bin"
ExecStart=/var/www/meta_app_factory/backend/venv/bin/uvicorn server:app --host 127.0.0.1 --port 5030
Restart=always

[Install]
WantedBy=multi-user.target
